# crypto-analyzer/systemd/smart-trader-v2.service
[Unit]
Description=Smart Trader Service V2
After=network.target mysql.service

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/crypto-analyzer
Environment="USE_V2_MODE=true"
ExecStart=/usr/bin/python3 /path/to/crypto-analyzer/smart_trader_service_v2.py
Restart=always
RestartSec=10

# 日志配置
StandardOutput=append:/path/to/crypto-analyzer/logs/v2_service.log
StandardError=append:/path/to/crypto-analyzer/logs/v2_service_error.log

[Install]
WantedBy=multi-user.target
