# crypto-analyzer/deploy/systemd/crypto-strategy-live.service
[Unit]
Description=Crypto Analyzer - 顶部反转做空策略 (strategy_live.py, ex-test3)
After=network-online.target mysql.service crypto-app-main.service
Wants=network-online.target
Requires=crypto-app-main.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/crypto-analyzer
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/crypto-analyzer/venv/bin/python strategy_live.py
Restart=on-failure
RestartSec=15
StandardOutput=append:/opt/crypto-analyzer/logs/strategy_live_systemd.log
StandardError=append:/opt/crypto-analyzer/logs/strategy_live_systemd.log

MemoryMax=1G
TasksMax=256

[Install]
WantedBy=multi-user.target
