# crypto-analyzer/deploy/systemd/crypto-strategy-bigmid.service
[Unit]
Description=Crypto Analyzer - Gemini AI 决策策略 (strategy_bigmid.py, ex-test3)
After=network-online.target mysql.service crypto-app-main.service
Wants=network-online.target
Requires=crypto-app-main.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/crypto-analyzer
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/crypto-analyzer/venv/bin/python strategy_bigmid.py
Restart=on-failure
RestartSec=30
StandardOutput=append:/opt/crypto-analyzer/logs/strategy_bigmid_systemd.log
StandardError=append:/opt/crypto-analyzer/logs/strategy_bigmid_systemd.log

MemoryMax=1G
TasksMax=256

[Install]
WantedBy=multi-user.target
