# crypto-analyzer/deploy/systemd/crypto-coin-futures.service
[Unit]
Description=Crypto Analyzer - 币本位主策略 (coin_futures_trader_service.py)
After=network-online.target mysql.service crypto-app-main.service
Wants=network-online.target
Requires=crypto-app-main.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/test2/crypto-analyzer
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/test2/crypto-analyzer/venv/bin/python coin_futures_trader_service.py
Restart=on-failure
RestartSec=15
StandardOutput=append:/home/test2/crypto-analyzer/logs/coin_futures_systemd.log
StandardError=append:/home/test2/crypto-analyzer/logs/coin_futures_systemd.log

MemoryMax=1G
TasksMax=512

[Install]
WantedBy=multi-user.target
