# crypto-analyzer/deploy/systemd/crypto-app-main.service
[Unit]
Description=Crypto Analyzer - FastAPI Web Server (app/main.py, port 9020)
Documentation=file:///home/test2/crypto-analyzer/design/MIGRATION_test3_to_test2.md
After=network-online.target mysql.service
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/test2/crypto-analyzer
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/test2/crypto-analyzer/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 9020
Restart=on-failure
RestartSec=10
# 优雅关闭：lifespan 会停后台任务；若 DB/API 卡住，最多等 120s 再 SIGKILL
TimeoutStopSec=120
StandardOutput=append:/home/test2/crypto-analyzer/logs/main_systemd.log
StandardError=append:/home/test2/crypto-analyzer/logs/main_systemd.log

# 资源限制 (防止内存泄漏拖死服务器)
MemoryMax=2G
TasksMax=512

[Install]
WantedBy=multi-user.target
