# crypto-analyzer/deploy/crypto-smart-trader.service
[Unit]
Description=Crypto Analyzer SmartTrader Service
After=network.target mysql.service mariadb.service crypto-fastapi.service
Wants=mysql.service mariadb.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/test2/crypto-analyzer
ExecStart=/home/test2/crypto-analyzer/venv/bin/python3 smart_trader_service.py
Restart=always
RestartSec=10
StandardOutput=append:/home/test2/crypto-analyzer/logs/smart_trader_stdout.log
StandardError=append:/home/test2/crypto-analyzer/logs/smart_trader_stderr.log

[Install]
WantedBy=multi-user.target
